# kolla-ansible/ansible/roles/cgroup_config/files/cgroup.slice
[Unit]
Description=Slice definition generated from the CEE configuration. Do not edit!
Documentation=man:systemd.special(7)
Before=slices.target

[Install]
WantedBy=slices.target
