# kolla-ansible/ansible/roles/cgroup_config/files/container.slice
[Unit]
Description=Entry point for controlling resources allocated to containerized services in CEE. Do not edit!
Documentation=man:systemd.special(7)
Before=slices.target

[Install]
WantedBy=slices.target
