# kolla-ansible/ansible/roles/cmha/templates/uptime.socket.j2
[Unit]
Description=Uptime Socket

[Socket]
ListenStream={{ agent_port_number }}
Accept=yes

[Install]
WantedBy=sockets.target