# kolla-ansible/ansible/roles/cmha/files/uptime@.service
[Unit]
Description=Uptime Service

[Service]
ExecStart=/bin/cat /proc/uptime
StandardOutput=socket
